# systemd/agentexw.service
[Unit]
Description=AgentExW Autonomous AI Agent
After=network.target whatsapp-svc.service email-svc.service calendar-svc.service
Wants=whatsapp-svc.service email-svc.service calendar-svc.service

[Service]
Type=simple
User=opc
Group=opc
Environment=ANTHROPIC_API_KEY=
EnvironmentFile=-/etc/execwall/env
ExecStart=/usr/local/bin/agentexw run --interval 300
Restart=always
RestartSec=10
WatchdogSec=120

# Systemd watchdog integration
NotifyAccess=main

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/execwall
PrivateTmp=true

[Install]
WantedBy=multi-user.target
