# systemd/agentexw-watcher.path
[Unit]
Description=Watch for AgentExW trigger events
After=agentexw.service
BindsTo=agentexw.service

[Path]
# Watch for new WhatsApp messages
PathChanged=/var/lib/execwall/events/whatsapp
# Watch for new emails
PathChanged=/var/lib/execwall/events/email
# Watch for fund alerts
PathChanged=/var/lib/execwall/events/fund_alerts
# Watch for manual trigger
PathExists=/var/lib/execwall/events/wake

[Install]
WantedBy=multi-user.target
