# systemd/calendar-svc.service
[Unit]
Description=Calendar Event Service for AgentExW
After=network.target

[Service]
Type=simple
User=opc
Group=opc
ExecStart=/usr/local/bin/calendar-svc
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/execwall
PrivateTmp=true

[Install]
WantedBy=multi-user.target
